# tasks/proactive-calendar-launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Template launchd job for the Proactive Calendar Conflicts one-shot
  (tasks/proactive-calendar.md).
  Copy to ~/Library/LaunchAgents/com.rachel.proactive-calendar.plist, replace
  __REPO_PATH__ with the absolute path to this repo checkout, then:
    launchctl load ~/Library/LaunchAgents/com.rachel.proactive-calendar.plist

  Runs bin/rachel headlessly with stdin closed, so the one-shot turn
  completes and the process exits on its own — no KeepAlive needed, each
  StartCalendarInterval entry fires a fresh run (same shape as
  tasks/inbox-brief-launchd.plist).

  Cadence: 4x/day during Gary's waking hours (Ireland time) — 08:00, 11:00,
  14:00, 17:00, matching the proactive layer's calendar_oneshot_hours
  {8, 11, 14, 17}. The proactive sweep (proactive/sweep.ts) also spawns
  wake-time catch-up runs for missed hours; when both fire, push.ts dedup
  absorbs the overlap (same conflict, same state => one ping).

  RACHEL_ALLOWED_TOOLS narrows the headless agent to the minimum set the
  task needs — Read, Write, Bash, and the Google Calendar MCP tools. The
  env var can only REMOVE tools from rachel.ts's default list, never add
  (injection hardening in proactive/allowedTools.ts).

  EnvironmentVariables/PATH is required because launchd jobs run with no
  user PATH, so tsx's `#!/usr/bin/env node` shebang can't resolve node
  (same requirement as bridge/launchd.plist).
-->
<plist version="1.0">
<dict>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>RACHEL_ALLOWED_TOOLS</key>
    <string>Read,Write,Bash,mcp__claude_ai_Google_Calendar__*</string>
  </dict>

  <key>Label</key>
  <string>com.rachel.proactive-calendar</string>

  <key>ProgramArguments</key>
  <array>
    <string>__REPO_PATH__/bin/rachel</string>
    <string>Read tasks/proactive-calendar.md and follow it.</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__REPO_PATH__</string>

  <key>StandardInPath</key>
  <string>/dev/null</string>

  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key>
      <integer>8</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>11</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>14</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>17</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
  </array>

  <key>StandardOutPath</key>
  <string>__REPO_PATH__/.rachel/proactive-calendar.log</string>

  <key>StandardErrorPath</key>
  <string>__REPO_PATH__/.rachel/proactive-calendar.log</string>
</dict>
</plist>
