# bridge/launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Template launchd job for the Telegram bridge — the single getUpdates
  consumer for the configured bot token. Copy to
  ~/Library/LaunchAgents/com.rachel.telegram-bridge.plist, replace
  __REPO_PATH__ with the absolute path to this repo checkout, then:
    launchctl load ~/Library/LaunchAgents/com.rachel.telegram-bridge.plist

  KeepAlive restarts the process on crash or on the loud exit(1) the bridge
  performs when it detects a second getUpdates consumer (HTTP 409).
  Logs are redacted by the bridge itself (bridge/api.ts's redact()) before
  ever reaching stdout/stderr, so these on-disk log files never carry the
  bot token even though launchd persists them.
  EnvironmentVariables/PATH is required because launchd jobs run with no
  user PATH, so tsx's `#!/usr/bin/env node` shebang can't resolve node.
  EnvironmentVariables/RACHEL_SESSION_FILE persists the Telegram session
  across a bridge restart (rachel.ts's RACHEL_SESSION_FILE seam) so a
  restart doesn't wipe the conversation thread. Set ONLY here — the CLI and
  the headless one-shots must never set this, since two processes resuming
  the same session id would fork the SDK's transcript chain.
-->
<plist version="1.0">
<dict>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>RACHEL_SESSION_FILE</key>
    <string>__REPO_PATH__/.rachel/bridge-session.json</string>
  </dict>

  <key>Label</key>
  <string>com.rachel.telegram-bridge</string>

  <key>ProgramArguments</key>
  <array>
    <string>__REPO_PATH__/node_modules/.bin/tsx</string>
    <string>__REPO_PATH__/bridge/telegram-bridge.ts</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__REPO_PATH__</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <!-- Belt-and-suspenders: even if the bridge exits on a 409, launchd won't
       restart it faster than 60s — giving Telegram time to release its
       getUpdates lock (~30-60s). The bridge's own backoff handles the race,
       but this prevents a future regression from re-creating the crash loop. -->
  <key>ThrottleInterval</key>
  <integer>60</integer>

  <key>StandardOutPath</key>
  <string>__REPO_PATH__/.rachel/telegram-bridge.log</string>

  <key>StandardErrorPath</key>
  <string>__REPO_PATH__/.rachel/telegram-bridge.log</string>
</dict>
</plist>
